# gnome/.config/autostart/irssi.desktop
[Desktop Entry]
Exec=/usr/bin/irssi
Name=irssi
NoDisplay=true
Terminal=true
Type=Application
